# templates/systemd/app.service.j2
# Systemd Service File (T-073)
# Generated by AXLE OS for application process management.

[Unit]
Description={{ app_name }} — managed by AXLE OS
After=network.target{% if database_service is defined %} {{ database_service }}.service{% endif %}

Wants=network-online.target

[Service]
Type=simple
User={{ user | default('www-data', true) }}
Group={{ group | default('www-data', true) }}
WorkingDirectory=/var/www/{{ app_name }}
ExecStart={{ start_command }}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app_name }}

# Environment
{% for env_var in env_vars | default([]) %}
Environment="{{ env_var }}"
{% endfor %}
{% if env_file is defined %}
EnvironmentFile={{ env_file }}
{% endif %}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/{{ app_name }}
{% if data_dir is defined %}
ReadWritePaths={{ data_dir }}
{% endif %}

[Install]
WantedBy=multi-user.target
